# deb/DEBIAN/postinst
if [ "$CLEAN_CADDYFILE" = "True" ]; then
    /usr/bin/echo "" > /etc/caddy/Caddyfile
fi

/bin/bash /var/lib/libre-workspace/portal/install.sh 

/usr/bin/systemctl daemon-reload 
/usr/bin/systemctl enable libre-workspace-portal.service
